# private/community_files/kyritz/files/etc/init.d/meshkitpostinstall
#!/bin/sh /etc/rc.common
# Freifunk Leipzig post run
# $Id: freifunk-leipzig  2014-02-21 23:51:08Z ufo $

START=99

start() {


chmod o+r /etc/config/olsrd
chmod o+r /etc/config/network
}

