# private/files/wizard.init
#!/bin/sh /etc/rc.common
# This script is part of Meshkit / Meshwizard
# It used to setup the Router after the first boot,
# after that it is disabled

START=99
boot() {
	logger -t wizard.notice \"the meshwizard is starting his work\"
	/usr/bin/meshwizard/wizard.sh
	/etc/init.d/wizard disable
	cp -r /rom/etc/config/meshkit /etc/config/
	logger -t wizard.notice \"The meshwizard has finished\"
}
